# _configurations/gunicorn.service
#save to     /etc/systemd/system/gunicorn.service

[Unit]
Description=gunicorn daemon
Requires=gunicorn.socket
After=network.target

[Service]
User=sammy
Group=www-data
WorkingDirectory=/home/ubuntu/encarta/encarta2k18/
ExecStart=/home/ubuntu/encarta/myprojectenv/bin/gunicorn \
          --access-logfile - \
          --workers 3 \
          --bind unix:/run/gunicorn.sock \
          encarts.wsgi:application

[Install]
WantedBy=multi-user.target